# conf/ibuntu_sddm.service
[Unit]
Description=ibuntu_sddmrenew
After=network.target

[Service]
ExecStart=/etc/init.d/ibuntu_sddm.sh

[Install]
WantedBy=default.target

